# docker/php-fpm/conf.d/transaction.conf
[program:transaction]
process_name=%(process_num)s
command=/usr/local/bin/php /var/www/artisan queue:work --queue=transaction --sleep=1 --tries=5 --timeout=24
autostart=true
autorestart=true
user=www
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/transaction.log
